# docker/vault_supervisord.ini
[program:consul]
priority=1
user=consul
environment=HOME="/etc/consul.d/",USER="consul"
command=/usr/local/bin/consul agent -config-file=/etc/consul.d/config.json -pid-file=/etc/consul.d/consul.pid
stderr_logfile = /tmp/results/consul.err.log
stdout_logfile = /tmp/results/consul.log

[program:vault-0.9.6-dev]
priority=2
user=vault-0.9.6-dev
directory=/etc/vault.d-0.9.6-dev/
environment=HOME="/etc/vault.d-0.9.6-dev/",USER="vault-0.9.6-dev"
command=/usr/local/bin/vault.0.9.6 server -dev -dev-listen-address=127.0.0.1:9200 -dev-root-token-id=${VAULT_V1_DEV_TOKEN}
stderr_logfile = /tmp/results/vault.0.9.6-dev.err.log
stdout_logfile = /tmp/results/vault.0.9.6-dev.log

[program:vault-1.5.4-dev]
priority=2
user=vault-1.5.4-dev
directory=/etc/vault.d-1.5.4-dev/
environment=HOME="/etc/vault.d-1.5.4-dev/",USER="vault-1.5.4-dev"
command=/usr/local/bin/vault.1.5.4 server -dev -dev-listen-address=127.0.0.1:9400 -dev-root-token-id=${VAULT_V2_DEV_TOKEN}
stderr_logfile = /tmp/results/vault.1.5.4-dev.err.log
stdout_logfile = /tmp/results/vault.1.5.4-dev.log

[program:vault-0.9.6-prod]
priority=3
user=vault-0.9.6-prod
environment=HOME="/etc/vault.d-0.9.6-prod/",USER="vault-0.9.6-prod"
command=/usr/local/bin/vault.0.9.6 server -config=/etc/vault.d-0.9.6-prod/config.hcl
stderr_logfile = /tmp/results/vault.0.9.6-prod.err.log
stdout_logfile = /tmp/results/vault.0.9.6-prod.log

[program:vault_1.5.4-prod]
priority=3
user=vault-1.5.4-prod
environment=HOME="/etc/vault.d-1.5.4-prod/",USER="vault-1.5.4-prod"
command=/usr/local/bin/vault.1.5.4 server -config=/etc/vault.d-1.5.4-prod/config.hcl
stderr_logfile = /tmp/results/vault.1.5.4-prod.err.log
stdout_logfile = /tmp/results/vault.1.5.4-prod.log
