# docker/vault_supervisord.ini
[program:consul]
priority=1
user=consul
environment=HOME="/etc/consul.d/",USER="consul"
command=/usr/local/bin/consul agent -config-file=/etc/consul.d/config.json -pid-file=/etc/consul.d/consul.pid
stderr_logfile = /tmp/results/consul.err.log
stdout_logfile = /tmp/results/consul.log

[program:vault-v1-dev]
priority=2
user=vault-v1-dev
directory=/etc/vault.d-v1-dev/
environment=HOME="/etc/vault.d-v1-dev/",USER="vault-v1-dev"
command=/usr/local/bin/vault.v1 server -dev -dev-listen-address=127.0.0.1:9200 -dev-root-token-id=devtoken1
stderr_logfile = /tmp/results/vault.v1-dev.err.log
stdout_logfile = /tmp/results/vault.v1-dev.log

[program:vault-v2-dev]
priority=2
user=vault-v2-dev
directory=/etc/vault.d-v2-dev/
environment=HOME="/etc/vault.d-v2-dev/",USER="vault-v2-dev"
command=/usr/local/bin/vault.v2 server -dev -dev-listen-address=127.0.0.1:9400 -dev-root-token-id=devtoken2
stderr_logfile = /tmp/results/vault.v2-dev.err.log
stdout_logfile = /tmp/results/vault.v2-dev.log

[program:vault-v1-prod]
priority=3
user=vault-v1-prod
environment=HOME="/etc/vault.d-v1-prod/",USER="vault-v1-prod"
command=/usr/local/bin/vault.v1 server -config=/etc/vault.d-v1-prod/config.hcl
stderr_logfile = /tmp/results/vault.v1-prod.err.log
stdout_logfile = /tmp/results/vault.v1-prod.log

[program:vault_v2-prod]
priority=3
user=vault-v2-prod
environment=HOME="/etc/vault.d-v2-prod/",USER="vault-v2-prod"
command=/usr/local/bin/vault.v2 server -config=/etc/vault.d-v2-prod/config.hcl
stderr_logfile = /tmp/results/vault.v2-prod.err.log
stdout_logfile = /tmp/results/vault.v2-prod.log
